# slideshow/scripts/slideshow-stream.service.example
[Unit]
Description=Pepites du jour - livestream YouTube en boucle
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/USER/slideshow
ExecStart=/home/USER/slideshow/scripts/stream_to_youtube.sh rtmp://a.rtmp.youtube.com/live2 TA_CLE_STREAM
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
